# deploy/systemd/spravka@.service
# Template unit — one instance per app:
#   systemctl enable --now spravka@public spravka@yurist spravka@admin spravka@rahbar spravka@qr
# Install to: /etc/systemd/system/spravka@.service
#
# All five. nginx routes bright.qrsystem.uz at :5000, so leaving spravka@qr out of this line
# means that subdomain answers 502 while everything else looks healthy.
#
# Each app's .env must carry its own PORT (qr 5000, public 5100, yurist 5101, admin 5102,
# rahbar 5103) — systemd expands ${PORT} below from the EnvironmentFile.

[Unit]
Description=Spravka web-%i
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=spravka
WorkingDirectory=/opt/spravka/apps/web-%i
EnvironmentFile=/opt/spravka/apps/web-%i/.env
Environment=NODE_ENV=production

# Signing renders the maʼlumotnoma with headless Chromium. Puppeteer caches its browser under
# $HOME by default, which ProtectHome= below makes unreachable — so the launch would fail at the
# first Imzolash, long after a deploy that looked fine. Keep the cache inside the app tree.
Environment=PUPPETEER_CACHE_DIR=/opt/spravka/.cache/puppeteer

# Not `npm run start`: that script hardcodes -H 0.0.0.0 for local dev, and CLI
# flags beat env vars — so in production every app would also answer on the
# public IP directly, over plain http, bypassing nginx and TLS entirely.
ExecStart=/usr/bin/npx next start -H 127.0.0.1 -p ${PORT}

Restart=always
RestartSec=5

# The app reads its own tree and writes only to .next, /tmp, and the document store.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/spravka/apps/web-%i/.next
# Issued maʼlumotnoma files (CERT_STORAGE_DIR). rahbar freezes one at each signing; public
# re-freezes documents signed before that existed. Leading '-' so yurist/admin, which never
# write here, still start if the directory is absent.
ReadWritePaths=-/var/lib/spravka/storage
# web-qr writes uploads and generated QR PNGs into its own public/ at runtime. ProtectSystem=strict
# makes everything else read-only, so without these the app starts fine and then fails on the first
# upload with EROFS. '-' because only web-qr has these directories.
ReadWritePaths=-/opt/spravka/apps/web-qr/public/uploads
ReadWritePaths=-/opt/spravka/apps/web-qr/public/qr

StandardOutput=journal
StandardError=journal
SyslogIdentifier=spravka-%i

[Install]
WantedBy=multi-user.target
